# recipes-revolution/revolution/service/revolution.service
[Unit]
Description=Revolution Autostart

[Service]
ExecStart=/usr/bin/revolution

[Install]
Alias=revolution.service
